# addon/supervisord.d/mix.queue.conf
[program:mix_queue]
process_name=%(process_num)02d
command=php /var/www/html/mix/artisan horizon
autostart=true
autorestart=true
user=root
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/mix/storage/logs/queue.log